# build/ipk/packages_arm/SmartThings-EE/usr/LXC/SmartThings-EE/etc/init.d/smartthing_detector
#!/bin/sh /etc/rc.common
START=3
STOP=3
start(){
	echo xx > /usr/xx
        service_start /usr/local/bin/smartthing_detector 
}

stop(){
    service_stop /usr/local/bin/smartthing_detector
}
