# build/ipk/packages_arm/smartthing_detector/etc/init.d/smartthing_detector
#!/bin/sh /etc/rc.common

start(){
        service_start /usr/local/bin/smartthing_detector 
}

stop(){
    service_stop /usr/local/bin/smartthing_detector
}
